# lang/supervisor/pheanstalk.ini
;官方消息发布事件消费者 
[program:pheanstalk]
command=/bin/php72 /data/releases/current/console/index.php pheanstalk/Run 
process_name=%(program_name)s_%(process_num)02d
numprocs=2
autostart=true
autorestart=true
stopsignal=TERM
stopwaitsecs=3
user=nobody
redirect_stderr=false
stdout_logfile=/data/log/daemon/daemon-stdout-%(program_name)s_%(process_num)02d.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
stderr_logfile=/data/log/daemon/daemon-stderr-%(program_name)s_%(process_num)02d.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=10
